# infra/containerbot/task-board-backup.service
[Unit]
Description=Back up Task Board state and source to the NAS
After=docker.service network-online.target

[Service]
Type=oneshot
User=serviceuser
WorkingDirectory=/srv/task-board/checkout
ExecStart=/usr/bin/flock --wait 600 /var/lib/task-board/backup.lock /bin/sh /srv/task-board/checkout/infra/containerbot/backup.sh
TimeoutStartSec=30min
